# system_services/init/caller.service
[Unit]
Description=HERMES Gateway caller daemon
After=nginx.service
Requires=nginx.service

[Service]
Type=exec
ExecStart=/usr/bin/caller.sh
KillSignal=SIGTERM
IgnoreSIGPIPE=no
Restart=always
TimeoutSec=infinity

[Install]
WantedBy=multi-user.target
