# system_services/init/modem.service
[Unit]
Description=HERMES open modem (mercury) daemon

[Service]
Type=exec
User=root
Group=root
<<<<<<< HEAD
ExecStart=/usr/bin/mercury -v
=======
ExecStart=/usr/bin/mercury
>>>>>>> a516f03 (systemd: simplify modem.service)
KillSignal=SIGTERM
IgnoreSIGPIPE=no
Restart=always

[Install]
WantedBy=multi-user.target
